# cmake-build-debug/polycubed.service
[Unit]
Description=polycube service

[Service]
Type=simple
ExecStart=/usr/local/bin/polycubed

[Install]
WantedBy=multi-user.target
